# agents/fusion-ai/deploy/linux/fusion-ai.service
[Unit]
Description=FusionAI — Multi-Model Teneo Agent (Gemini + Groq + Claude + GPT)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=teneo
WorkingDirectory=/opt/fusion-ai
EnvironmentFile=/opt/fusion-ai/.env
ExecStart=/opt/fusion-ai/deploy/linux/start-agent.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=fusion-ai

[Install]
WantedBy=multi-user.target
